# backend/comfyui-api-service/comfyui-api.service
[Unit]
Description=ComfyUI API Service
After=network.target comfyui.service
Requires=comfyui.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/comfyui_api_service
Environment="PATH=/home/ubuntu/ComfyUI/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="S3_BUCKET=your-bucket-name"
Environment="AWS_REGION=us-east-1"
ExecStart=/home/ubuntu/ComfyUI/venv/bin/python /home/ubuntu/comfyui_api_service/api_service.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/comfyui-api.log
StandardError=append:/var/log/comfyui-api-error.log

[Install]
WantedBy=multi-user.target
